# src/common/salt/states/500-docker-gcf-cron/supervisor.conf
[program:cron]
command=cron-logs
autostart=true
priority=500
stdout_logfile=syslog
stdout_logfile_maxbytes=0
stderr_logfile=syslog
stderr_logfile_maxbytes=0
